# etc/rc.local
#!/bin/bash
# File: /etc/rc.local
# Execute: chmod +x /etc/rc.local
caddy file-server --root /root/pi-extras --browse &
exit 0
